# services/kblight.service
[Unit]
Description=Turn on scrolllock based keyboar light

[Service]
ExecStart=/usr/sbin/kblight

[Install]
WantedBy=multi.user.target
